# minecraft-backup/crontab
0 * * * * cd /data && /usr/local/bin/minecraft_autobackup.sh

